# supervisord.conf
[supervisord]
nodaemon=true

[program:events]
command=java -jar -Xms1G -Xmx2G /opt/zoo-event-1.0.0-SNAPSHOT/zoo-event.jar /opt/zoo-event-1.0.0-SNAPSHOT/conf.edn
