# supervisord.conf
[supervisord]
nodaemon=true

[program:events]
command=java -jar -Xms1G -Xmx2G /opt/zoo-events/zoo-events-0.3.0.jar /opt/zoo-events/conf.edn
