# supervisord.conf
[supervisord]
nodaemon=true

[program:events]
command=java -jar -Xms1G -Xmx2G /opt/zoo-event/zoo-event.jar /opt/zoo-event/conf/conf.edn
